# deployment/systemd/ai-news-tracker.timer
[Unit]
Description=AI News Tracker Pipeline Timer
Documentation=https://github.com/yourusername/ai-news-tracker
Requires=ai-news-tracker.service

[Timer]
# Run every hour
OnCalendar=hourly
# Run immediately if missed (e.g., system was off)
Persistent=true
# Add random delay to avoid system load spikes (0-60 seconds)
RandomizedDelaySec=60

[Install]
WantedBy=timers.target

